# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --daemon
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true

[program:phpfpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true